# docs/runit/sv/aasd/run
#!/bin/sh
[ -r conf ] && . ./conf
exec 2>&1
sleep 30
exec /opt/AASd/aasd.py ${OPTS}
